# deploy.timer
[Unit]
Description=Update Signin Web App deployment Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target